# packages/hostScriptsServices/hostServices/docker-upgrade.service
[Unit]
Description=Update docker engine service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/src/dappnode/DNCORE/scripts/host/docker_upgrade.sh --upgrade

[Install]
WantedBy=multi-user.target
