# packages/hostScriptsServices/hostServices/update-upgrade-host.service
[Unit]
Description=Update and upgrade DAppNode host

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/src/dappnode/DNCORE/scripts/host/update_upgrade.sh

[Install]
WantedBy=multi-user.target
