# packages/hostScriptsServices/hostTimers/check-docker-network.timer
[Unit]
Description=Timer for checking Docker network status
After=network.target

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=check-docker-network.service

[Install]
WantedBy=timers.target